# config/systemd/user/omarchy-chromium-shutdown.service
[Unit]
Description=Graceful Chromium shutdown for Omarchy
DefaultDependencies=no
Before=shutdown.target reboot.target halt.target wayland-wm@hyprland.desktop.service
Conflicts=shutdown.target reboot.target halt.target
PartOf=wayland-session@hyprland.desktop.target

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/bin/true
ExecStop=/bin/bash -c 'echo "Gracefully shutting down Chromium processes..."; if command -v hyprctl >/dev/null 2>&1; then echo "Closing Chromium windows via Hyprland..."; hyprctl clients -j | jq -r ".[] | select(.class | test(\"chromium|chrome\"; \"i\")) | .address" | while read addr; do hyprctl dispatch closewindow "address:$addr" 2>/dev/null || true; done; sleep 3; fi; echo "Sending SIGTERM to remaining processes..."; for proc in chromium chromium-browser chrome google-chrome; do pkill -TERM "$proc" 2>/dev/null && echo "Sent SIGTERM to $proc" || true; done; sleep 5; echo "Force killing any remaining processes..."; for proc in chromium chromium-browser chrome google-chrome; do pkill -KILL "$proc" 2>/dev/null && echo "Force killed $proc" || true; done; echo "Chromium shutdown complete"'
TimeoutStopSec=20

[Install]
WantedBy=wayland-session@hyprland.desktop.target