# config/systemd/user/omarchy-chromium-shutdown.service
[Unit]
Description=Graceful Chromium shutdown for Omarchy
DefaultDependencies=no
Before=shutdown.target reboot.target halt.target wayland-wm@hyprland.desktop.service
Conflicts=shutdown.target reboot.target halt.target
PartOf=wayland-session@hyprland.desktop.target

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/bin/true
ExecStop=/home/arch/.local/bin/omarchy-chromium-shutdown-helper
TimeoutStopSec=20

[Install]
WantedBy=wayland-session@hyprland.desktop.target