# config/systemd/user/omarchy-chromium-shutdown.service
[Unit]
Description=Graceful Chromium shutdown for Omarchy
DefaultDependencies=no
Before=shutdown.target reboot.target halt.target

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/bin/true
ExecStop=/bin/bash -c 'for proc in chromium chromium-browser chrome google-chrome; do pkill -TERM "$proc" 2>/dev/null || true; done; sleep 3; for proc in chromium chromium-browser chrome google-chrome; do pkill -KILL "$proc" 2>/dev/null || true; done'
TimeoutStopSec=10

[Install]
WantedBy=default.target